# poly-bot-wrap.service
[Unit]
Description=5m Poly Bot - Polymarket Auto Wrap (USDC.e -> pUSD)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/root/5m-poly-bot
Environment="PATH=/root/5m-poly-bot/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/usr/bin/npm run wrap
StandardOutput=append:/root/5m-poly-bot/wrap-worker.log
StandardError=append:/root/5m-poly-bot/wrap-worker.log

