# poly-bot-live.service
[Unit]
Description=5m Poly Bot — Crypto Trading Bot (LIVE)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/5m-poly-bot
EnvironmentFile=/root/5m-poly-bot/.env
Environment="PATH=/root/5m-poly-bot/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTHONUNBUFFERED=1"
ExecStart=/root/5m-poly-bot/.venv/bin/python -u /root/5m-poly-bot/crypto_bot.py --live
Restart=always
RestartSec=10
StandardOutput=append:/root/5m-poly-bot/bot-live.log
StandardError=append:/root/5m-poly-bot/bot-live.log

[Install]
WantedBy=multi-user.target
